# scripts/com.user.stock-analysis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.stock-analysis</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/liyuhang/.local/bin/python3.11</string>
        <string>/Users/liyuhang/Documents/stock/main.py</string>
        <string>start</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>8000</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/liyuhang/Documents/stock</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/liyuhang/Documents/stock/data/logs/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/liyuhang/Documents/stock/data/logs/launchd.err.log</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
</dict>
</plist>